# systemd/origincards.service
[Unit]
Description=Origin Cards Image Generator Service
After=network.target

[Service]
Type=simple
User=errantson
Group=errantson
WorkingDirectory=/srv/reverie.house/services/origincards
ExecStart=/usr/local/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=origincards

# Security hardening
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
ReadWritePaths=/srv/reverie.house/services/origincards

[Install]
WantedBy=multi-user.target
