# systemd/reverie-mapperwatch.service
[Unit]
Description=Reverie House - Mapperwatch (Origins Quest Monitor)
Documentation=https://github.com/reverie-house/reverie
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=errantson
WorkingDirectory=/srv/reverie.house
ExecStart=/usr/bin/python3 core/mapperwatch.py --verbose
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reverie-mapperwatch
Environment="PYTHONUNBUFFERED=1"

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
