# ops/contributions.service
[Unit]
Description=Reverie Contribution Score Update Service
After=network.target

[Service]
Type=oneshot
ExecStart=/srv/reverie.house/ops/manage_contributions.sh quiet
WorkingDirectory=/srv/reverie.house
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
